# tools/mac/local.ghandwa-db-backup-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.ghandwa-db-backup-pull</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/g_/Library/Scripts/pull-ghandwa-db-backups.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/g_/Library/Logs/ghandwa-db-backup/launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/g_/Library/Logs/ghandwa-db-backup/launchd.log</string>
</dict>
</plist>
